# com.tabelog.checker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tabelog.checker</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Volumes/mitomo_SSD/tabelog_checker/venv/bin/python3</string>
        <string>/Volumes/mitomo_SSD/tabelog_checker/check_scheduled.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Volumes/mitomo_SSD/tabelog_checker</string>
    <key>RunAtLoad</key>
    <false/>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Volumes/mitomo_SSD/tabelog_checker/logs/check.log</string>
    <key>StandardErrorPath</key>
    <string>/Volumes/mitomo_SSD/tabelog_checker/logs/check_error.log</string>
</dict>
</plist>
